# systemd/forge-dashboard.service
[Unit]
Description=Forge Prime Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=FORGE_PRIME_DIR
ExecStart=/usr/bin/python3 -m uvicorn dashboard.app:app --host FORGE_DASHBOARD_BIND --port 8099
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
